# battery-alert.service
[Unit]
Description=Battery Alert Monitor (GPIO 25)
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/battery-alert.py
Restart=always
RestartSec=5
User=root
Group=root

# Allow the service to shutdown the system
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=battery-alert

[Install]
WantedBy=multi-user.target
